# root/etc/services.d/hyperhdr/run
#!/usr/bin/with-contenv bash

# set umask
UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec \
        hyperhdr -v \
	-u=/config
