# root/etc/services.d/hyperhdr/run
exec \
	s6-setuidgid abc /usr/bin/hyperiond \
	-u=/config
